# templates/default/twistd_plugin.conf.erb
description                     "Twistd service: <%= @service_name %>"
author                          "Ying Li <cyli@twistedmatrix.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

<% if @user != "root" %>
exec su - <%= @user %> -c '<%= "authbind --deep " if @authbinded %>twistd -n \
<% else %>
exec twistd -n \
<% end %>
    --pidfile <%= @pidfile %> \
<% if @logfile %>
    --logfile <%= @logfile %> \
<% end %>
    <%= @twistd_command -%>
<% if @args.length > 0 %>
<%=     " \\" if @args.length > 0 %>
<% else %>
<%=     "'" if @user != "root" %>
<% end %>
<% @args.each_with_index do |argline, index| %>
    <%= argline -%>
<%      if @args.length - 1 == index %>
<%=         "'" if @user != "root" -%>
<%=         "\n" -%>
<%      else %>
<%=         " \\\n" -%>
<%      end %>
<% end %>
