# templates/default/twistd_plugin.conf.erb
description                     "Twistd service: <%= @service_name %>"
author                          "Ying Li <cyli@twistedmatrix.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

<% if @user != "root" %>
exec su - <%= @user %> -c '<%= "authbind --deep " if @authbinded %>twistd -n \
<% else %>
exec twistd -n \
<% end %>
<% if @log_dir.nil? %>
    --pidfile "/tmp/<%= @service_name %>.twistd.pid"
<% end %>
<% if @log_dir %>
    --pidfile <%= @log_dir %>/twistd.pid \
    --logfile <%= @log_dir %>/twistd.log \
<% end %>
    <%= @twistd_command %> \
<% @args.each do |argline| %>
    <%= argline %> \
<% end %><%= "'" if @user != "root" %>
